# family-finance-admin-bot.service
[Unit]
Description=Family Finance Admin Panel Telegram Bot
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/family-finance-bot
Environment="PATH=/opt/family-finance-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/family-finance-bot/.env
ExecStart=/opt/family-finance-bot/venv/bin/python /opt/family-finance-bot/admin_bot.py
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=family-finance-admin-bot

# Security settings (same as main bot)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/family-finance-bot/logs /opt/family-finance-bot/family_finance.db
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target

